# airflow.conf
[program:airflow_web]
command=/usr/local/bin/airflow webserver -p 8080

[program:airflow_scheduler]
command=/usr/local/bin/airflow scheduler